# day.cron
#!/bin/bash

cd /var/www/missloopy

source venv/bin/activate

sh fetch-bannedips.sh
sh fetch-scammers.sh

python email-newmembers.py
#python purge-bounced.py

make
